# deploy/systemd/worker.service
[Unit]
Description=Uptime Monitoring Worker
After=network.target
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/uptime-monitoring
EnvironmentFile=/etc/uptime-monitoring.env
ExecStart=/usr/bin/env python -m app.workers.runner
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
